# static/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
echo "Starting NGinx"
exec /usr/sbin/nginx -c /config/nginx/nginx.conf
