# static/etc/services.d/irssi/run
#!/usr/bin/with-contenv bash

umask 0002
export TERM=xterm

chown -R bytesized:bytesized /config/irssi

exec screen -D -m -S irssi s6-setuidgid bytesized irssi --home /config/irssi
